# seamail.conf
[program:seamail]
command=/opt/seamail/seamail.py %(process_num)02d
process_name=%(program_name)s-%(process_num)s
numprocs=1
numprocs_start=0
autostart=true
autorestart=true
startsecs=10
